# deployment/roles/backend_build/files/scheduler_service_docker/supervisord.conf
[supervisord]
nodaemon=true
loglevel=info
[unix_http_server]
file=/run/supervisord.sock

[supervisorctl]
serverurl=unix:///run/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:scheduler-service]
command=/home/allu/scheduler-service/run_scheduler-service.sh
user=allu
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true

[eventlistener:stdout]
command = supervisor_stdout
buffer_size = 100
events = PROCESS_LOG
result_handler = supervisor_stdout:event_handler

