# deployment/roles/frontend/files/frontend_docker/supervisord.conf
[supervisord]
nodaemon=true

[program:apache]
command = /home/allu/apache/apache_start.sh
user = root
autostart = true
